# ffremote.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>KeepAlive</key>
  <true/>
  <key>Label</key>
  <string>ffremote</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>ASPNETCORE_ENVIRONMENT</key>
    <string>Production</string>
    <key>ASPNETCORE_URLS</key>
    <string>http://localhost:@PORT@</string>
    <key>DOTNET_BUNDLE_EXTRACT_BASE_DIR</key>
    <string>/var/tmp/.net</string>
    <key>DOTNET_PRINT_TELEMETRY_MESSAGE</key>
    <string>false</string>
  </dict>
  <key>ProgramArguments</key>
  <array>
    <string>@BINDIR@/@ASSEMBLY_NAME@</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/var/log/@ASSEMBLY_NAME@.log</string>
  <key>UserName</key>
  <string>@USER@</string>
  <key>WorkingDirectory</key>
  <string>@SYSCONFDIR@/@ASSEMBLY_NAME@</string>
</dict>
</plist>
